# deploy/systemd/scum-bot.service
[Unit]
Description=SCUM Platform Bot
After=network.target

[Service]
Type=simple
User=scum
WorkingDirectory=/opt/scum-platform
Environment=NODE_ENV=production
Environment=PERSIST_REQUIRE_DB=true
Environment=PERSIST_LEGACY_SNAPSHOTS=false
EnvironmentFile=/opt/scum-platform/.env
ExecStart=/usr/bin/node src/bot.js
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
