# deploy/systemd/scum-watcher.service
[Unit]
Description=SCUM Platform Watcher
After=network.target

[Service]
Type=simple
User=scum
WorkingDirectory=/opt/scum-platform
Environment=NODE_ENV=production
EnvironmentFile=/opt/scum-platform/.env
ExecStart=/usr/bin/node src/services/scumLogWatcherRuntime.js
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
